# velican/conf/systemd.service
[Unit]
Description=Watch %i for changes and regenerate

[Service]
ExecStart={{BINARY}} %i
