# velican/conf/systemd.timer
[Timer]
OnActiveSec=75days
RandomizedDelaySec=120